# systemd/net-mgr-relay.service
[Unit]
Description=net-mgr-relay — replicate from peer net-mgr instances
Documentation=https://github.com/kev-cam/net-mgr
After=network-online.target net-mgr.service mariadb.service mysql.service
Wants=network-online.target
PartOf=net-mgr.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/net-mgr-relay --foreground
# always (not on-failure): a clean exit 0 must still be restarted — a transient
# empty roster once left this dead for days. The daemon also no longer exits 0
# while it can still watch/serve (see sbin/net-mgr-relay idle-watch mode).
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=net-mgr-relay

[Install]
WantedBy=multi-user.target
